# snapraid-scrub.service
[Unit]
Description=SnapRAID Scrub Service
[Service]
Type=oneshot
ExecStart=/usr/bin/snapraid scrub
[Install]
WantedBy=multi-user.target