# snapraid-scrub.timer
[Unit]
Description=Run SnapRAID scrub service nightly at 1AM

[Timer]
OnCalendar=*-*-* 1:00:00
Persistent=true

[Install]
WantedBy=timers.target